# system/startup/infinite-candle.service
[Unit]
Description=Infinite Candle LED Service
After=network.target sound.target

[Service]
Type=simple
# 1. MUST be root to access hardware PWM (/dev/mem)
User=root
Group=root

# 2. Set working directory so relative paths work
WorkingDirectory=/home/pi/infinite-candle/lights-control

# 3. CRITICAL: Use the python binary INSIDE the .venv folder
#    We also use '-m src.main' to resolve the import paths correctly
ExecStart=/home/pi/infinite-candle/lights-control/.venv/bin/python -m src.main

# 4. Environment variables (optional, but good for safety)
Environment="PATH=/home/pi/infinite-candle/lights-control/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target